# wn/image/profile.d/cudalib.sh
export PATH=${PATH}:/usr/local/nvidia/bin:/usr/local/cuda/bin:/usr/local/bin:/usr/bin:/bin
export LD_LIBRARY_PATH=${LD_LIBRARY_PATH}:/usr/lib64:/usr/local/cuda/lib64/:/usr/local/nvidia/lib64/
